# certbot/crontab
45 3,15 * * * root certbot certonly --standalone -d $DOMAINS > /proc/1/fd/1 2>/proc/1/fd/2
